# files/.config/systemd/user/clementine.service
[Unit]
Description=Clementine

[Service]
ExecStart=/usr/bin/clementine

[Install]
WantedBy=default.target
