# files/.config/systemd/user/compton.service
[Unit]
Description=compton

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=default.target
