# files/.config/systemd/user/slack.service
[Unit]
Description=Slack

[Service]
ExecStart=/usr/bin/flatpak run com.slack.Slack

[Install]
WantedBy=default.target
