# files/.config/systemd/user/slack.service
[Unit]
Description=Slack

[Service]
ExecStart=/usr/bin/slack

[Install]
WantedBy=default.target
