# files/.config/systemd/user/thunderbird.service
[Unit]
Description=Thunderbird

[Service]
ExecStart=/usr/bin/flatpak run org.mozilla.Thunderbird

[Install]
WantedBy=default.target
