# files/.config/systemd/user/thunderbird.service
[Unit]
Description=Thunderbird

[Service]
ExecStart=/usr/bin/thunderbird

[Install]
WantedBy=default.target
