# files/.config/systemd/user/xfce4-notifyd.service
[Unit]
Description=xfce4 notify daemon

[Service]
ExecStart=/usr/lib64/xfce4/notifyd/xfce4-notifyd

[Install]
WantedBy=default.target
