# files_root/etc/systemd/system/slow.slice
[Unit]
Description=Slice used to run slow and stupid stuff.
Before=slices.target

[Slice]
CPUAccounting=true
CPUQuota=5%

# Defaults to 100, so this means 10%.
IOAccounting=true
IOWeight=5
