# files_root/etc/systemd/system/sentinelone.service.d/override.conf
[Service]
Slice=slow.slice
TimeoutStopSec=1s